# kazunya.service
[Unit]
Description=Kazunya Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kazunya
WorkingDirectory=/opt/kazunya
ExecStart=/usr/bin/java -jar /opt/kazunya/Kazunya-1.2.3.jar
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/kazunya

[Install]
WantedBy=multi-user.target